# Library/LaunchAgents/com.user.no-apple-music.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- This script stops Apple Music from auto-launching when media keys are pressed -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://apple.com">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.no-apple-music</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/launchctl</string>
        <string>unload</string>
        <string>-w</string>
        <string>/System/Library/LaunchAgents/com.apple.rcd.plist</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <false/>
</dict>
</plist>

